# htpc-session.service
[Unit]
Description=HTPC kiosk session (cage + launcher)
After=systemd-user-sessions.service network-online.target
Wants=network-online.target
Conflicts=getty@tty1.service

[Service]
User=htpc
PAMName=login
TTYPath=/dev/tty1
StandardInput=tty
StandardOutput=journal
UtmpIdentifier=tty1
UtmpMode=user
Environment=XDG_RUNTIME_DIR=/run/user/1000
ExecStart=/usr/bin/cage -d -- /opt/htpc/system/start-session.sh
Restart=always
RestartSec=2

[Install]
WantedBy=graphical.target
